# systemd/tars-deploy.service
[Unit]
Description=TARS auto-deploy (pull main and restart services)

[Service]
Type=oneshot
WorkingDirectory=%h/tars-voice-agent
ExecStart=%h/tars-voice-agent/scripts/deploy.sh
